# Library/LaunchAgents/com.codex.obsidian-sync.plist.tmpl
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.codex.obsidian-sync</string>
  <key>ProgramArguments</key>
  <array>
    <string>/bin/bash</string>
    <string>{{ .chezmoi.homeDir }}/.codex/hooks/watch-and-save.sh</string>
  </array>
  <key>RunAtLoad</key>
  <true/>
  <key>KeepAlive</key>
  <true/>
  <key>StandardOutPath</key>
  <string>{{ .chezmoi.homeDir }}/.codex/logs/obsidian-sync.log</string>
  <key>StandardErrorPath</key>
  <string>{{ .chezmoi.homeDir }}/.codex/logs/obsidian-sync.err</string>
</dict>
</plist>
